# deployments/linux-host/systemd/postgresql-16.service.d/onfailure.conf
# Drop-in override — adds Telegram alerting to the package-shipped
# postgresql-16 unit without modifying it directly. Install at
# /etc/systemd/system/postgresql-16.service.d/onfailure.conf then
# `systemctl daemon-reload`.
[Unit]
OnFailure=notify-telegram@%n.service
