# deployments/linux-host/systemd/mixchamb.service
[Unit]
Description=Mixchamb Phoenix app
After=network-online.target postgresql-16.service cloudflared.service
Wants=network-online.target
Requires=postgresql-16.service

[Service]
Type=exec
User=mixchamb
Group=mixchamb
EnvironmentFile=/etc/mixchamb.env
WorkingDirectory=/opt/mixchamb
ExecStart=/opt/mixchamb/bin/server
ExecStop=/opt/mixchamb/bin/mixchamb stop
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
OnFailure=notify-telegram@%n.service

# Hardening — replaces most of what a container would give you
# for isolation on a single-app box.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
ReadWritePaths=/var/log/mixchamb

[Install]
WantedBy=multi-user.target
